# recipes-extended/logrotate/logrotate/clean-logrotate.cron
#!/bin/sh

test -f /var/lib/logrotate.status && rm /var/lib/logrotate.status
